# theme/luci-theme-atmaterial/root/etc/uci-defaults/30_luci-theme-atmaterial
#!/bin/sh
uci batch <<-EOF
	set luci.themes.ATMaterial_Brown=/luci-static/atmaterial_Brown
	commit luci
EOF
exit 0
